# hickory-dns/files/etc/init.d/hickory-dns
#!/bin/sh /etc/rc.common

USE_PROCD=1
START=18
STOP=51

PROG=/usr/bin/hickory-dns
CONF=/etc/hickory-dns/forwarder.toml
TSIG_KEY=/etc/hickory-dns/update.key

generate_tsig_key() {
    [ -f "$TSIG_KEY" ] && return 0
    dd if=/dev/urandom bs=32 count=1 of="$TSIG_KEY" 2>/dev/null
    chown nobody:nogroup "$TSIG_KEY" 2>/dev/null || true
    chmod 400 "$TSIG_KEY"
}

start_service() {
    generate_tsig_key
    # Remove stale journal so hickory-dns always starts from the zone file.
    # ipv6-neigh re-registers all dynamic entries on startup, so the journal
    # has no durable value and may contain a corrupt state after unclean shutdown.
    rm -f /var/lib/hickory-dns/lan.jrnl
    procd_open_instance hickory-dns
    procd_set_param command $PROG -c $CONF -q
    procd_set_param stdout 1
    procd_set_param stderr 1
    procd_set_param respawn "${respawn_threshold:-3600}" "${respawn_timeout:-5}" "${respawn_retry:-5}"
    procd_close_instance hickory-dns
}

start() {
    rc_procd start_service "$@"
    # Journal was wiped above; restart ipv6-neigh so it re-registers all
    # dynamic DNS entries from the current neighbour table.
    /etc/init.d/ipv6-neigh restart
}


service_triggers() {
    procd_add_reload_trigger "dhcp" "network" "system"
    procd_add_reload_interface_trigger "wan_6"
}
