# luci-app-pppoe-user/root/etc/init.d/pppoe-user
#!/bin/sh /etc/rc.common
# Copyright (C) 2024-2026 OpenWrt.org
# Licensed under GPL v2

START=90
STOP=10

USE_PROCD=1
PROG=/usr/sbin/pppoe-server
UBUS_HELPER=/usr/share/ucode/pppoe-ubus.uc
UCODE_BIN=/usr/bin/ucode

SPEED_LIMIT_SWITCH=1
CHAP_SECRETS=/etc/ppp/chap-secrets
BACKUP_DIR=/var/pppoe-user/backup
LOG_DIR=/var/pppoe-user/log
SESSION_DIR=/var/etc/pppoe-user/session

# 初始化目录
init_dirs() {
    mkdir -p "$BACKUP_DIR"
    mkdir -p "$LOG_DIR"
    mkdir -p "$LOG_DIR/userinfo"
    mkdir -p "$SESSION_DIR"
    mkdir -p /var/lock
}

# 使用 ucode 进行日期比较和状态检查 (更高效且不易出错)
check_expiration_and_build_secrets() {
    local current_date
    current_date=$(date +"%Y-%m-%d")
    
    # 清空并重建 chap-secrets
    : > "$CHAP_SECRETS"
    
    # 遍历 UCI 配置
    config_load pppoe-user
    
    config_foreach setup_user "$current_date"
}

# 回调函数：处理每个用户配置
# 参数 1: cfg_name, 参数 2: current_date
setup_user() {
    local cfg="$1"
    local today="$2"
    local enabled username servicename password ipaddr expires
    
    config_get_bool enabled "$cfg" enabled 0
    config_get username "$cfg" username
    config_get servicename "$cfg" servicename
    config_get password "$cfg" password
    config_get ipaddr "$cfg" ipaddr
    config_get expires "$cfg" expires
    
    # 默认 servicename 为 *
    [ -z "$servicename" ] && servicename="*"
    
    # 检查过期
    if [ "$enabled" -eq 1 ] && [ -n "$expires" ] && [ "$expires" = "$today" ]; then
        logger -t pppoe-user "User $username expired today ($today). Disabling."
        uci set pppoe-user."$cfg".enabled=0
        uci commit pppoe-user
        enabled=0
    fi
    
    # 如果启用且未过期，写入 chap-secrets
    if [ "$enabled" -eq 1 ]; then
        # 格式: user service secret IP
        echo "$username $servicename $password $ipaddr" >> "$CHAP_SECRETS"
    fi
}

# 管理 Crontab 任务
manage_cron() {
    local cron_cmd_start="59 6 * * * /etc/init.d/pppoe-user start"
    local cron_cmd_backup="59 23 * * * cp -f /etc/config/pppoe-user $BACKUP_DIR/pppoe-user.\$(date +\%Y\%m\%d).bak"
    
    # 移除旧的任务（防止重复）
    crontab -u root -l 2>/dev/null | grep -v "/etc/init.d/pppoe-user start" | grep -v "cp -f /etc/config/pppoe-user" | crontab -u root -
    
    # 添加新任务
    (crontab -u root -l 2>/dev/null; echo "$cron_cmd_start"; echo "$cron_cmd_backup") | sort -u | crontab -u root -
}

# 配置 nftables 限速表
setup_nft() {
    # 检查表是否存在，存在则先删除以保证幂等性
    nft list table inet pppoe >/dev/null 2>&1 && nft delete table inet pppoe
    
    nft add table inet pppoe
    nft add chain inet pppoe upload { type filter hook prerouting priority filter\; policy accept\; }
    nft add chain inet pppoe download { type filter hook postrouting priority filter\; policy accept\; }
    
    # 更新防火墙区域配置以包含 ppp 接口
    local zone_device
    zone_device=$(uci get firewall.@zone[0].device 2>/dev/null)
    if [ -z "$(echo "$zone_device" | grep 'ppp+')" ]; then
        uci add_list firewall.@zone[0].device='ppp+'
        uci commit firewall
        # 注意：这里不重启 firewall，以免断开现有连接，nft 规则即时生效
    fi
}

cleanup_nft() {
    nft list table inet pppoe >/dev/null 2>&1 && nft delete table inet pppoe
}

# 启动 ucode ubus 辅助服务
start_ubus_helper() {
    if [ -x "$UCODE_BIN" ] && [ -f "$UBUS_HELPER" ]; then
        # 检查是否已在运行
        if ! pgrep -f "ucode.*pppoe-ubus.uc" >/dev/null; then
            logger -t pppoe-user "Starting ubus helper service"
            $UCODE_BIN $UBUS_HELPER &
            disown
        fi
    else
        logger -t pppoe-user "Warning: ucode helper not found or not executable. LuCI kill function may fail."
    fi
}

stop_ubus_helper() {
    pkill -f "ucode.*pppoe-ubus.uc"
}

start_service() {
    init_dirs
    
    logger -t pppoe-user "Starting PPPoE User Management Service"
    
    # 备份旧的 chap-secrets
    [ -f "$CHAP_SECRETS" ] && cp -f "$CHAP_SECRETS" /etc/ppp/lasttime-chap-secrets
    
    # 生成新的 chap-secrets 并检查过期
    check_expiration_and_build_secrets
    
    # 设置权限
    chmod 600 "$CHAP_SECRETS"
    
    # 配置定时任务
    manage_cron
    
    # 配置限速
    if [ "$SPEED_LIMIT_SWITCH" = "1" ]; then
        setup_nft
    else
        cleanup_nft
    fi
    
    # 启动 ubus 辅助进程 (供 LuCI JS 调用)
    start_ubus_helper
    
    # 如果是作为服务运行，不需要保持前台，因为主要工作是生成配置和规则
    # 真正的 PPPoE 服务由 pppoe-server 守护进程管理
}

stop_service() {
    logger -t pppoe-user "Stopping PPPoE User Management Service"
    
    stop_ubus_helper
    cleanup_nft
    
    # 可选：清除 chap-secrets (通常不建议，因为重启后需要，但在 stop 时清理可确保安全)
    # : > "$CHAP_SECRETS" 
}

reload_service() {
    stop_service
    start_service
}

# 兼容旧版 procd 的 stop 行为
stop() {
    stop_service
}
