# cgroupfs-mount/files/cgroupfs-mount.init
#!/bin/sh /etc/rc.common

START=01

boot() {
	cgroupfs-mount
}
