# luci-app-passwall/root/etc/uci-defaults/luci-passwall.rej
--- luci-app-passwall/root/etc/uci-defaults/luci-passwall
+++ luci-app-passwall/root/etc/uci-defaults/luci-passwall
@@ -1,5 +1,8 @@
 #!/bin/sh
 
+grep -q ip-api.com /usr/share/passwall/rules/proxy_host ||
+	sed -i '$a ip-api.com' /usr/share/passwall/rules/proxy_host
+
 uci -q batch <<-EOF >/dev/null
 	set dhcp.@dnsmasq[0].localuse=1
 	commit dhcp
